# systemd/llm-chat-frontend.service
[Unit]
Description=LLM Chat System Frontend Service
After=network.target llm-chat-backend.service

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/llm-chat-system/frontend
Environment="PATH=/usr/bin:/usr/local/bin:/usr/local/bin/node"
Environment="NODE_ENV=production"
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10

# 日志配置
StandardOutput=append:/path/to/llm-chat-system/logs/frontend.log
StandardError=append:/path/to/llm-chat-system/logs/frontend-error.log

[Install]
WantedBy=multi-user.target
